# scripts/telegram-diary-summarize.service
[Unit]
Description=Telegram Diary LLM Summarize
After=network.target

[Service]
Type=oneshot
User=muffin
WorkingDirectory=/home/muffin/prod/telegram_diary
ExecStart=/home/muffin/prod/telegram_diary/scripts/summarize.sh
EnvironmentFile=/home/muffin/prod/telegram_diary/.env
StandardOutput=journal
StandardError=journal

# セキュリティ強化
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/muffin/prod/telegram_diary
ReadWritePaths=/home/muffin/.cache/uv
ReadWritePaths=/home/muffin/dev/notes/200_Personal/Diary
