# scripts/telegram-diary.service
[Unit]
Description=Telegram Diary Polling Service
After=network.target

[Service]
Type=simple
User=muffin
WorkingDirectory=/home/muffin/prod/telegram_diary
ExecStart=/home/muffin/.local/bin/uv run python -m src.main
EnvironmentFile=/home/muffin/prod/telegram_diary/.env
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# セキュリティ強化
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
# 作業ディレクトリへの書き込みを許可（logs/, messages/, daily/, state.json）
ReadWritePaths=/home/muffin/prod/telegram_diary
# uv がキャッシュへ書き込めるようにする
ReadWritePaths=%h/.cache/uv

[Install]
WantedBy=multi-user.target
